# deploy/systemd/bithumb-bot-healthcheck.service
[Unit]
Description=Bithumb bot healthcheck
After=bithumb-bot.service

[Service]
Type=oneshot
WorkingDirectory=@BITHUMB_BOT_ROOT@
Environment=BITHUMB_ENV_FILE=@BITHUMB_ENV_FILE_LIVE@
ExecStart=/usr/bin/env uv run python @BITHUMB_BOT_ROOT@/scripts/healthcheck.py
