# deploy/systemd/bithumb-bot-healthcheck.service
[Unit]
Description=Bithumb bot healthcheck
After=bithumb-bot.service

[Service]
Type=oneshot
WorkingDirectory=/workspace/bithumb-bot
EnvironmentFile=-/etc/bithumb-bot/bithumb-bot.env
ExecStart=/usr/bin/env bash -lc 'uv run python scripts/healthcheck.py'
