# deploy/systemd/bithumb-bot-backup.service
[Unit]
Description=Backup Bithumb bot SQLite database

[Service]
Type=oneshot
WorkingDirectory=@BITHUMB_BOT_ROOT@
Environment=BITHUMB_ENV_FILE=@BITHUMB_ENV_FILE_LIVE@
ExecStart=/usr/bin/env bash @BITHUMB_BOT_ROOT@/scripts/backup_sqlite.sh
