# deploy/systemd/bithumb-bot-paper.service
[Unit]
Description=Bithumb bot trading loop (paper)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=@BITHUMB_BOT_ROOT@
Environment="MODE=paper" "BITHUMB_ENV_FILE=@BITHUMB_ENV_FILE_PAPER@" "PYTHONUNBUFFERED=1" "ENV_ROOT=@BITHUMB_ENV_ROOT@" "RUN_ROOT=@BITHUMB_RUN_ROOT@" "DATA_ROOT=@BITHUMB_DATA_ROOT@" "LOG_ROOT=@BITHUMB_LOG_ROOT@" "BACKUP_ROOT=@BITHUMB_BACKUP_ROOT@"
ExecStart=/usr/bin/env bash -lc 'exec @BITHUMB_UV_BIN@ run python -u -m bithumb_bot run --mode paper --short ${SMA_SHORT:-7} --long ${SMA_LONG:-30}'
Restart=always
RestartSec=5
SyslogIdentifier=bithumb-bot-paper
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
